# systemd/container-ed363dfa21375e837cbb865354d5d3ba0df85165024a9f76240cb944987c793f.service
# container-ed363dfa21375e837cbb865354d5d3ba0df85165024a9f76240cb944987c793f.service
# autogenerated by Podman 3.2.0-rc1
# Wed May 12 18:52:42 BST 2021

[Unit]
Description=Podman container-ed363dfa21375e837cbb865354d5d3ba0df85165024a9f76240cb944987c793f.service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
After=network-online.target
RequiresMountsFor=/var/home/pod_server/.local/share/containers/storage /run/user/1001/containers
BindsTo=pod-d4bb53335436c00538732543fc60df32d7225dac022187f5d55060d2416fe684.service
After=pod-d4bb53335436c00538732543fc60df32d7225dac022187f5d55060d2416fe684.service

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
TimeoutStopSec=70
ExecStart=/usr/bin/podman start ed363dfa21375e837cbb865354d5d3ba0df85165024a9f76240cb944987c793f
ExecStop=/usr/bin/podman stop -t 10 ed363dfa21375e837cbb865354d5d3ba0df85165024a9f76240cb944987c793f
ExecStopPost=/usr/bin/podman stop -t 10 ed363dfa21375e837cbb865354d5d3ba0df85165024a9f76240cb944987c793f
PIDFile=/run/user/1001/containers/overlay-containers/ed363dfa21375e837cbb865354d5d3ba0df85165024a9f76240cb944987c793f/userdata/conmon.pid
Type=forking

[Install]
WantedBy=multi-user.target default.target
