# systemd/container-776d637846cc6dfaeb6dbd106807d09bb1c516e8c5d1d9f30e76f4823ae492c0.service
# container-776d637846cc6dfaeb6dbd106807d09bb1c516e8c5d1d9f30e76f4823ae492c0.service
# autogenerated by Podman 3.2.0-rc1
# Wed May 12 18:52:42 BST 2021

[Unit]
Description=Podman container-776d637846cc6dfaeb6dbd106807d09bb1c516e8c5d1d9f30e76f4823ae492c0.service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
After=network-online.target
RequiresMountsFor=/var/home/pod_server/.local/share/containers/storage /run/user/1001/containers
BindsTo=pod-d4bb53335436c00538732543fc60df32d7225dac022187f5d55060d2416fe684.service
After=pod-d4bb53335436c00538732543fc60df32d7225dac022187f5d55060d2416fe684.service

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
TimeoutStopSec=70
ExecStart=/usr/bin/podman start 776d637846cc6dfaeb6dbd106807d09bb1c516e8c5d1d9f30e76f4823ae492c0
ExecStop=/usr/bin/podman stop -t 10 776d637846cc6dfaeb6dbd106807d09bb1c516e8c5d1d9f30e76f4823ae492c0
ExecStopPost=/usr/bin/podman stop -t 10 776d637846cc6dfaeb6dbd106807d09bb1c516e8c5d1d9f30e76f4823ae492c0
PIDFile=/run/user/1001/containers/overlay-containers/776d637846cc6dfaeb6dbd106807d09bb1c516e8c5d1d9f30e76f4823ae492c0/userdata/conmon.pid
Type=forking

[Install]
WantedBy=multi-user.target default.target
