# systemd/container-a6aed7be5a8a87eef7d01f1bda1f350d36b88c93151f79e7440533ef44086819.service
# container-a6aed7be5a8a87eef7d01f1bda1f350d36b88c93151f79e7440533ef44086819.service
# autogenerated by Podman 3.2.0-rc1
# Wed May 12 18:52:42 BST 2021

[Unit]
Description=Podman container-a6aed7be5a8a87eef7d01f1bda1f350d36b88c93151f79e7440533ef44086819.service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
After=network-online.target
RequiresMountsFor=/var/home/pod_server/.local/share/containers/storage /run/user/1001/containers
BindsTo=pod-d4bb53335436c00538732543fc60df32d7225dac022187f5d55060d2416fe684.service
After=pod-d4bb53335436c00538732543fc60df32d7225dac022187f5d55060d2416fe684.service

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
TimeoutStopSec=70
ExecStart=/usr/bin/podman start a6aed7be5a8a87eef7d01f1bda1f350d36b88c93151f79e7440533ef44086819
ExecStop=/usr/bin/podman stop -t 10 a6aed7be5a8a87eef7d01f1bda1f350d36b88c93151f79e7440533ef44086819
ExecStopPost=/usr/bin/podman stop -t 10 a6aed7be5a8a87eef7d01f1bda1f350d36b88c93151f79e7440533ef44086819
PIDFile=/run/user/1001/containers/overlay-containers/a6aed7be5a8a87eef7d01f1bda1f350d36b88c93151f79e7440533ef44086819/userdata/conmon.pid
Type=forking

[Install]
WantedBy=multi-user.target default.target
