# systemd/container-006d5a42e7ec7a52aade0e33488ae4b6632445a3517cd3e59a847273aad324f4.service
# container-006d5a42e7ec7a52aade0e33488ae4b6632445a3517cd3e59a847273aad324f4.service
# autogenerated by Podman 3.2.0-rc1
# Wed May 12 18:52:42 BST 2021

[Unit]
Description=Podman container-006d5a42e7ec7a52aade0e33488ae4b6632445a3517cd3e59a847273aad324f4.service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
After=network-online.target
RequiresMountsFor=/var/home/pod_server/.local/share/containers/storage /run/user/1001/containers
BindsTo=pod-d4bb53335436c00538732543fc60df32d7225dac022187f5d55060d2416fe684.service
After=pod-d4bb53335436c00538732543fc60df32d7225dac022187f5d55060d2416fe684.service

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
TimeoutStopSec=70
ExecStart=/usr/bin/podman start 006d5a42e7ec7a52aade0e33488ae4b6632445a3517cd3e59a847273aad324f4
ExecStop=/usr/bin/podman stop -t 10 006d5a42e7ec7a52aade0e33488ae4b6632445a3517cd3e59a847273aad324f4
ExecStopPost=/usr/bin/podman stop -t 10 006d5a42e7ec7a52aade0e33488ae4b6632445a3517cd3e59a847273aad324f4
PIDFile=/run/user/1001/containers/overlay-containers/006d5a42e7ec7a52aade0e33488ae4b6632445a3517cd3e59a847273aad324f4/userdata/conmon.pid
Type=forking

[Install]
WantedBy=multi-user.target default.target
