# systemd/container-9af3a90db94c5c10519d6ec01c72a9372fc07510ab4f3b7e38725381848214c6.service
# container-9af3a90db94c5c10519d6ec01c72a9372fc07510ab4f3b7e38725381848214c6.service
# autogenerated by Podman 3.2.0-rc1
# Wed May 12 18:52:42 BST 2021

[Unit]
Description=Podman container-9af3a90db94c5c10519d6ec01c72a9372fc07510ab4f3b7e38725381848214c6.service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
After=network-online.target
RequiresMountsFor=/var/home/pod_server/.local/share/containers/storage /run/user/1001/containers
BindsTo=pod-d4bb53335436c00538732543fc60df32d7225dac022187f5d55060d2416fe684.service
After=pod-d4bb53335436c00538732543fc60df32d7225dac022187f5d55060d2416fe684.service

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
TimeoutStopSec=70
ExecStart=/usr/bin/podman start 9af3a90db94c5c10519d6ec01c72a9372fc07510ab4f3b7e38725381848214c6
ExecStop=/usr/bin/podman stop -t 10 9af3a90db94c5c10519d6ec01c72a9372fc07510ab4f3b7e38725381848214c6
ExecStopPost=/usr/bin/podman stop -t 10 9af3a90db94c5c10519d6ec01c72a9372fc07510ab4f3b7e38725381848214c6
PIDFile=/run/user/1001/containers/overlay-containers/9af3a90db94c5c10519d6ec01c72a9372fc07510ab4f3b7e38725381848214c6/userdata/conmon.pid
Type=forking

[Install]
WantedBy=multi-user.target default.target
