# systemd/container-009bbe701dce61f14f9de8e375c088ccd7a3a3ce74e6320dedb6b8a517dbc606.service
# container-009bbe701dce61f14f9de8e375c088ccd7a3a3ce74e6320dedb6b8a517dbc606.service
# autogenerated by Podman 3.2.0-rc1
# Wed May 12 18:52:42 BST 2021

[Unit]
Description=Podman container-009bbe701dce61f14f9de8e375c088ccd7a3a3ce74e6320dedb6b8a517dbc606.service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
After=network-online.target
RequiresMountsFor=/var/home/pod_server/.local/share/containers/storage /run/user/1001/containers
BindsTo=pod-d4bb53335436c00538732543fc60df32d7225dac022187f5d55060d2416fe684.service
After=pod-d4bb53335436c00538732543fc60df32d7225dac022187f5d55060d2416fe684.service

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
TimeoutStopSec=70
ExecStart=/usr/bin/podman start 009bbe701dce61f14f9de8e375c088ccd7a3a3ce74e6320dedb6b8a517dbc606
ExecStop=/usr/bin/podman stop -t 10 009bbe701dce61f14f9de8e375c088ccd7a3a3ce74e6320dedb6b8a517dbc606
ExecStopPost=/usr/bin/podman stop -t 10 009bbe701dce61f14f9de8e375c088ccd7a3a3ce74e6320dedb6b8a517dbc606
PIDFile=/run/user/1001/containers/overlay-containers/009bbe701dce61f14f9de8e375c088ccd7a3a3ce74e6320dedb6b8a517dbc606/userdata/conmon.pid
Type=forking

[Install]
WantedBy=multi-user.target default.target
