# systemd/pod-d4bb53335436c00538732543fc60df32d7225dac022187f5d55060d2416fe684.service
# pod-d4bb53335436c00538732543fc60df32d7225dac022187f5d55060d2416fe684.service
# autogenerated by Podman 3.2.0-rc1
# Wed May 12 18:52:42 BST 2021

[Unit]
Description=Podman pod-d4bb53335436c00538732543fc60df32d7225dac022187f5d55060d2416fe684.service
Documentation=man:podman-generate-systemd(1)
Wants=network.target
After=network-online.target
RequiresMountsFor= 
Requires=container-006d5a42e7ec7a52aade0e33488ae4b6632445a3517cd3e59a847273aad324f4.service container-009bbe701dce61f14f9de8e375c088ccd7a3a3ce74e6320dedb6b8a517dbc606.service container-776d637846cc6dfaeb6dbd106807d09bb1c516e8c5d1d9f30e76f4823ae492c0.service container-9af3a90db94c5c10519d6ec01c72a9372fc07510ab4f3b7e38725381848214c6.service container-a6aed7be5a8a87eef7d01f1bda1f350d36b88c93151f79e7440533ef44086819.service container-ed363dfa21375e837cbb865354d5d3ba0df85165024a9f76240cb944987c793f.service
Before=container-006d5a42e7ec7a52aade0e33488ae4b6632445a3517cd3e59a847273aad324f4.service container-009bbe701dce61f14f9de8e375c088ccd7a3a3ce74e6320dedb6b8a517dbc606.service container-776d637846cc6dfaeb6dbd106807d09bb1c516e8c5d1d9f30e76f4823ae492c0.service container-9af3a90db94c5c10519d6ec01c72a9372fc07510ab4f3b7e38725381848214c6.service container-a6aed7be5a8a87eef7d01f1bda1f350d36b88c93151f79e7440533ef44086819.service container-ed363dfa21375e837cbb865354d5d3ba0df85165024a9f76240cb944987c793f.service

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
TimeoutStopSec=70
ExecStart=/usr/bin/podman start e67b5ac2cbea33ad9f672b78ea23a8c358e07463188ea59c3fd9736b40a948bf
ExecStop=/usr/bin/podman stop -t 10 e67b5ac2cbea33ad9f672b78ea23a8c358e07463188ea59c3fd9736b40a948bf
ExecStopPost=/usr/bin/podman stop -t 10 e67b5ac2cbea33ad9f672b78ea23a8c358e07463188ea59c3fd9736b40a948bf
PIDFile=/run/user/1001/containers/overlay-containers/e67b5ac2cbea33ad9f672b78ea23a8c358e07463188ea59c3fd9736b40a948bf/userdata/conmon.pid
Type=forking

[Install]
WantedBy=multi-user.target default.target
